# services/systemd/9beads-user.service
[Unit]
Description=9beads Task Server (User Service)
Documentation=https://github.com/lneely/9beads
After=default.target

[Service]
Type=simple
WorkingDirectory=%h

ExecStart=%h/bin/Beads fgstart

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
