# services/runit/run
#!/bin/sh
exec 2>&1

export HOME=${HOME:-/home/$(whoami)}
export PATH="${HOME}/bin:/usr/local/bin:/usr/bin:/bin"

exec chpst -u $(whoami) 9beads fgstart
